# scripts/com.wakwak.memex-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wakwak.memex-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wakwak/src/github.com/wakwak3125/memex/scripts/run-skill.sh</string>
        <string>/sync-all</string>
    </array>

    <!-- 08:00-23:00 毎時実行 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>3</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/wakwak/src/github.com/wakwak3125/memex/logs/launchd-sync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wakwak/src/github.com/wakwak3125/memex/logs/launchd-sync-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/wakwak/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/wakwak</string>
    </dict>
</dict>
</plist>
