# scripts/com.wakwak.memex-reflect-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wakwak.memex-reflect-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wakwak/src/github.com/wakwak3125/memex/scripts/run-skill.sh</string>
        <string>/reflect daily</string>
    </array>

    <!-- 平日 22:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>22</integer><key>Minute</key><integer>7</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/wakwak/src/github.com/wakwak3125/memex/logs/launchd-reflect-daily-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wakwak/src/github.com/wakwak3125/memex/logs/launchd-reflect-daily-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/wakwak/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/wakwak</string>
    </dict>
</dict>
</plist>
