# scripts/com.wakwak.memex-sync-meetings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wakwak.memex-sync-meetings</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wakwak/src/github.com/wakwak3125/memex/scripts/run-skill.sh</string>
        <string>/sync-meetings</string>
    </array>

    <!-- 毎日 18:00 実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/wakwak/src/github.com/wakwak3125/memex/logs/launchd-sync-meetings-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wakwak/src/github.com/wakwak3125/memex/logs/launchd-sync-meetings-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/wakwak/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/wakwak</string>
    </dict>
</dict>
</plist>
